# scripts/com.jiawen.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jiawen.morning-brief</string>
  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>__REPO_ROOT__/morning_brief.py</string>
    <string>--no-open</string>
    <string>--commit</string>
    <string>--push</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/.launchd.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
